# infra/yarrr-tech.service
[Unit]
Description=Yarrr.Tech FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/yarrr-tech
Environment="PYTHONPATH=/opt/yarrr-tech"
Environment="ALLOWED_ORIGIN=https://yarrr-node.com"
ExecStart=/opt/yarrr-tech/backend/.venv/bin/uvicorn backend.app:app --host 127.0.0.1 --port 8421 --workers 1 --log-level info
Restart=always
RestartSec=3
StandardOutput=append:/var/log/yarrr-tech.log
StandardError=append:/var/log/yarrr-tech.err.log

[Install]
WantedBy=multi-user.target
